# packaging/arc.service
[Unit]
Description=Arc - High-Performance Time-Series Data Warehouse
After=network.target
Documentation=https://github.com/Basekick-Labs/arc

[Service]
Type=simple
User=arc
Group=arc
WorkingDirectory=/opt/arc
ExecStart=/usr/local/bin/arc
Restart=on-failure
RestartSec=5s

# Environment
Environment="PATH=/opt/arc/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"

# Resource limits
LimitNOFILE=65536

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=arc

[Install]
WantedBy=multi-user.target
